# researchsite.service
[Unit]
Description=Research Site Flask Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/ResearchBuddy
Environment=PORT=80
Environment=PYTHONPATH=/ResearchBuddy
Environment=HF_HOME=/tmp/hf_cache
Environment=TRANSFORMERS_CACHE=/tmp/hf_cache
ExecStart=/ResearchBuddy/start.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target